# scripts/pixelpostman.service
[Unit]
Description=PixelPostman - QR interactive book + NFC travel postcards
After=network.target sound.target

[Service]
Type=simple
WorkingDirectory=@REPO_DIR@
ExecStart=@REPO_DIR@/venv/bin/python main.py
Restart=on-failure
RestartSec=3
User=@USER@
# Grant access to the scanner, NFC reader, HDMI-CEC and GPIO without root.
SupplementaryGroups=input i2c video gpio

[Install]
WantedBy=multi-user.target
